# launchd/local.kanata.builtin-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>local.kanata.builtin-watcher</string>
    <key>ProgramArguments</key>
    <array>
        <string>__HOME__/.config/kanata/nuphy_builtin_keyboard_watcher.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Library/Logs/local.kanata.builtin-watcher.log</string>
    <key>StandardErrorPath</key>
    <string>/Library/Logs/local.kanata.builtin-watcher.log</string>
</dict>
</plist>
